# system/services/jarvis-pipeline.service
[Unit]
Description=JARVIS Pipeline Engine Daemon
After=network.target jarvis-ws.service

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis-linux
ExecStart=/home/turbo/jarvis-linux/.venv/bin/python /home/turbo/jarvis-linux/scripts/pipeline_engine.py --daemon
Restart=always
RestartSec=30
Environment=TURBO=/home/turbo/jarvis-linux
Environment=PYTHONPATH=/home/turbo/jarvis-linux

[Install]
WantedBy=default.target
