# system/services/jarvis-telegram.service
[Unit]
Description=JARVIS Telegram Bot Bridge
After=network.target jarvis-proxy.service

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis-linux/canvas
ExecStart=/usr/bin/node /home/turbo/jarvis-linux/canvas/telegram-bot.js
Restart=always
RestartSec=10
Environment=TURBO=/home/turbo/jarvis-linux
Environment=TURBO_DIR=/home/turbo/jarvis-linux

[Install]
WantedBy=default.target
