# system/services/jarvis-voice.service
[Unit]
Description=JARVIS Voice Listener (Wake Word)
After=network.target jarvis-whisper.service jarvis-openclaw.service

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis-linux
Environment=PYTHONPATH=/home/turbo/jarvis-linux
Environment=OPENCLAW_GATEWAY_URL=http://127.0.0.1:18790/query
Environment=WHISPER_FLOW_HTTP_URL=http://127.0.0.1:9742/api/voice/transcribe_blob
Environment=JARVIS_TTS_SCRIPT=/home/turbo/jarvis-linux/scripts/jarvis-tts.sh
ExecStart=/home/turbo/jarvis-linux/.venv/bin/python3 /home/turbo/jarvis-linux/src/voice_engine.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
