# system/services/jarvis-whisper.service
[Unit]
Description=JARVIS WhisperFlow Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis-linux
Environment=PYTHONPATH=/home/turbo/jarvis-linux
ExecStart=/home/turbo/jarvis-linux/.venv/bin/python -m uvicorn python_ws.server:app --host 127.0.0.1 --port 9742
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
