# system/services/jarvis-mcp.service
[Unit]
Description=JARVIS Jarvis-mcp
After=network.target

[Service]
Type=simple
ExecStart=/home/turbo/jarvis-linux/.venv/bin/python -m src.mcp_server
WorkingDirectory=/home/turbo/jarvis-linux
Restart=always
RestartSec=5
Environment=PYTHONPATH=/home/turbo/jarvis-linux
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
