# system/services/jarvis-openclaw.service
[Unit]
Description=OpenClaw Gateway (LM Studio local)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/turbo
Environment=HOME=/home/turbo
Environment=OPENCLAW_STATE_DIR=/home/turbo/.openclaw
ExecStart=/usr/bin/openclaw gateway --port 18790
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
