# system/services/jarvis-backup.service
[Unit]
Description=JARVIS Backup Service

[Service]
Type=oneshot
ExecStart=/home/turbo/jarvis/scripts/jarvis-backup.sh
