# system/services/jarvis-proxy.service
[Unit]
Description=JARVIS Jarvis-proxy
After=network.target

[Service]
ExecStart=/home/turbo/jarvis-m1-ops/scripts/jarvis-boot.sh proxy
WorkingDirectory=/home/turbo/jarvis-m1-ops
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
