# system/services/jarvis-gemini-openai.service
[Unit]
Description=JARVIS Gemini OpenAI Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis
ExecStart=/usr/bin/node gemini-openai-proxy.js
Restart=always
RestartSec=5
Environment=TURBO=/home/turbo/jarvis

[Install]
WantedBy=default.target
