# system/services/jarvis-n8n.service
[Unit]
Description=JARVIS n8n Automation Engine
After=network.target

[Service]
Type=simple
Environment=N8N_PORT=5678
ExecStart=/home/turbo/.local/bin/n8n start
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
