# system/services/jarvis-health.timer
[Unit]
Description=Run JARVIS Cluster Health every 30min

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
