# system/services/jarvis-pipeline.service
[Unit]
Description=JARVIS Pipeline Engine Daemon
After=network.target jarvis-ws.service

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis
ExecStart=/home/turbo/jarvis/.venv/bin/python scripts/pipeline_engine.py --daemon
Restart=always
RestartSec=30
Environment=TURBO=/home/turbo/jarvis

[Install]
WantedBy=default.target
