# system/services/jarvis-proxy.service
[Unit]
Description=JARVIS Jarvis-proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis-linux/canvas
ExecStart=/usr/bin/node /home/turbo/jarvis-linux/canvas/direct-proxy.js
Restart=always
RestartSec=5
Environment=TURBO=/home/turbo/jarvis-linux
Environment=TURBO_DIR=/home/turbo/jarvis-linux
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
