# system/services/jarvis-watchdog.service
[Unit]
Description=JARVIS Master Watchdog
After=network.target jarvis-ws.service

[Service]
Type=simple
WorkingDirectory=/home/turbo/jarvis
ExecStart=/home/turbo/jarvis/.venv/bin/python -m scripts.openclaw_watchdog
Restart=always
RestartSec=60
Environment=TURBO=/home/turbo/jarvis
Environment=TURBO_DIR=/home/turbo/jarvis

[Install]
WantedBy=default.target
